# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/airbnb/bin/uwsgi --ini /srv/airbnb/.config/uwsgi/uwsgi-app.ini
